# templates/provisioner-cache.service.j2
[Unit]
Description=Provisioner Cache
After=docker.service
Requires=docker.service

[Service]
WorkingDirectory=/opt/provisioner-cache
ExecStart=/usr/bin/docker compose -f /opt/provisioner-cache/docker-compose.yaml up -d
ExecStop=/usr/bin/docker compose -f /opt/provisioner-cache/docker-compose.yaml down
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
